# nimdnsd.service
[Unit]
Description=NI mDNS Responder service

[Service]
Type=forking
ExecStart=/usr/bin/nimdnsResponder

[Install]
WantedBy=multi-user.target